# cron.hourly/twitter2mail
#!/bin/sh

export HOME=/home/danil
export RBENV_VERSION=2.2.2

su danil \
   --shell=/bin/bash \
   --command "cd ${HOME}/bin/thtm && \
              source ${HOME}/.bash_aliases && \
              bundle exec 'ruby $HOME/bin/thtm/thtm.rb'"

EXITVALUE=$?
if [ $EXITVALUE != 0 ]; then
    logger -t git "ALERT mirroring exited abnormally with [$EXITVALUE]"
fi

exit 0
