# systemd/system/git-public-mirror.service
[Unit]
Description="Git mirror of the public repositories"

[Service]
Type=oneshot
User=git
Group=git
ExecStart=/home/danil/bin/git-public-mirror
