# systemd/system/snap-telegram\x2ddesktop-2469.mount
[Unit]
Description=Mount unit for telegram-desktop, revision 2469
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/telegram-desktop_2469.snap
Where=/snap/telegram-desktop/2469
Type=squashfs
Options=nodev,ro,x-gdu.hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
