# init.d/swapfile
#!/sbin/openrc-run
# Copyright 1999-2016 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

depend()
{
	after localmount
}

start()
{
	ebegin "Activating swap file"
	swapon ${PATH_TO_SWAPFILE}
	eend $?
}

stop()
{
	ebegin "Deactivating swap file"
	swapoff ${PATH_TO_SWAPFILE}
	eend $?
}
